# config/systemd/sales-ops-celery.service
[Unit]
Description=Sales Ops Celery Worker
After=postgresql.service valkey.service network.target

[Service]
Type=simple
User=celery_user
Group=celery_user
WorkingDirectory=/opt/sales-ops/backend
ExecStart=/opt/sales-ops/venv/bin/celery -A config worker \
    --loglevel=info \
    --concurrency=4 \
    --max-tasks-per-child=1000 \
    --time-limit=300 \
    --soft-time-limit=240
Restart=always
RestartSec=10
MemoryMax=4G
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sales-ops-celery
EnvironmentFile=/opt/sales-ops/backend/.env

[Install]
WantedBy=multi-user.target
