# config/systemd/sales-ops-celery-beat.service
[Unit]
Description=Sales Ops Celery Beat (Scheduler)
After=postgresql.service valkey.service network.target

[Service]
Type=simple
User=celery_user
Group=celery_user
WorkingDirectory=/opt/sales-ops/backend
ExecStart=/opt/sales-ops/venv/bin/celery -A config beat \
    --loglevel=info \
    --pidfile=/var/run/sales-ops/celery-beat.pid
Restart=always
RestartSec=10
MemoryMax=512M
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sales-ops-celery-beat
EnvironmentFile=/opt/sales-ops/backend/.env

[Install]
WantedBy=multi-user.target
