# config/systemd/sales-ops-backend.service
[Unit]
Description=Sales Ops Backend (Uvicorn)
After=postgresql.service valkey.service network.target

[Service]
Type=notify
User=app_user
Group=app_user
WorkingDirectory=/opt/sales-ops/backend
ExecStart=/opt/sales-ops/venv/bin/uvicorn config.wsgi:application \
    --socket /var/run/sales-ops/app.sock \
    --workers 4 \
    --access-log - \
    --log-level info \
    --timeout-keep-alive 5
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
MemoryMax=2G
LimitNOFILE=65535
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sales-ops-backend
EnvironmentFile=/opt/sales-ops/backend/.env

[Install]
WantedBy=multi-user.target
